# files/var/spool/cron/crontabs/root
* * * * * /app/bin/speedtest-to-librato
